# scripts/piups-sensor.service
[Unit]
Description=PiUPS Breakout Garden IoT Sensor Streaming
Documentation=https://github.com/tspannhw
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=__PROJECT_DIR__
ExecStart=__VENV_DIR__/bin/python -m app.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=piups-sensor

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=LOG_LEVEL=INFO

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=__PROJECT_DIR__
PrivateTmp=true

# Resource limits for 2GB Pi
MemoryMax=256M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
